# config/99-wifi-hotplug.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0bda", RUN+="/etc/raspap/hostapd/servicestart.sh --interface br0 --seconds 3"
